# skills/claude-state-agent/systemd/claude-state-agent.service
[Unit]
Description=Local Claude state agent (HTTP — claude-fleet-starter)
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
# User=, WorkingDirectory=, ExecStart= are filled in by install.sh based on the
# user who ran the install. The template here uses placeholders so we don't ship
# /root hardcoded (the kit must run for any user, root or not).
User=__USER__
WorkingDirectory=__SKILL_DIR__
EnvironmentFile=-/etc/default/claude-state-agent
EnvironmentFile=-__INSTALL_DIR__/state-agent.env
ExecStartPre=/bin/bash -c 'until tailscale ip -4 2>/dev/null | grep -q .; do sleep 2; done'
ExecStart=__SKILL_DIR__/.venv/bin/python -m src.main
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
