# skills/claude-launcher/systemd/claude-launcher.service
[Unit]
Description=Auto-launch claude code in tmux at boot
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
Environment=HOME=/root
Environment=PATH=/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/root/skills/claude-launcher/scripts/launch.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
